# systemd-watchdog.service
[Unit]
Description=systemd watchdog service

[Service]
ExecStart=/usr/bin/dotnet /opt/daemons/systemd-watchdog/systemd-watchdog.dll
